# .config/autostart/make-cache.desktop
[Desktop Entry]
Exec=bash /home/chris/opt/make-cache.sh
Name=make-cache
StartupNotify=true
Terminal=false
Type=Application
